# deploy/systemd/agentifui.target
[Unit]
Description=AgentifUI application stack
Wants=agentifui-gateway.service agentifui-web.service
After=network.target

[Install]
WantedBy=multi-user.target
